# mealie-recipe-export.service.d/variables.conf
[Service]
Environment=MEALIE_URL=
Environment=MEALIE_API_TOKEN=
Environment=MEALIE_EXPORT_DIR=
