# mealie-recipe-import.path
[Unit]
Description=Watch for new or modified recipe files

[Path]
# Add PathModified directives to the provided *.conf file.
# Example,
# PathModified=/var/recipes/userA/*.json

[Install]
WantedBy=multi-user.target
