# mealie-recipe-export.timer
[Unit]
Description=Nightly export of Mealie recipes

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
